# deployment/ktlo-simulator.service
[Unit]
Description=KTLO Simulator Application
After=network.target

[Service]
Type=simple
User=manish
WorkingDirectory=/home/manish/ktlo_java_sim
ExecStart=/usr/bin/java -Xms512m -Xmx1024m -jar /home/manish/ktlo_java_sim/ktlo-simulator-1.0.0.jar
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

Environment="JAVA_HOME=/usr/lib/jvm/java-11-openjdk-amd64"

# Exit code 143 is a normal termination (SIGTERM)
SuccessExitStatus=143

# Timeout for stopping the service
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
